# daemon/deployer.service
[Unit]
Description=A deployment utility
Requires=polo.service
After=polo.service

[Service]

ExecStart=/opt/deployer/deployer/deployer.py &> /dev/null
ExecStop=/opt/deployer/daemon/deployerstop.py &> /dev/null

Type=simple
User=root
Group=root
PIDFile=/var/run/marcopolo/deployer.pid

[Install]
WantedBy=multi-user.target
