# daemon/old/systemd/marcoreceiverd.service
[Unit]
Description=A deployment utility
Requires=polo.service
After=polo.service

[Service]

Type=simple
User=root
Group=root
PIDFile=/var/run/marcoreceiverd.pid

[Install]
WantedBy=multi-user.target
